# stow/autostart/.config/autostart/vocalinux.desktop
[Desktop Entry]
Version=1.0
Type=Application
Name=Vocalinux
Exec=/home/pete/.local/bin/vocalinux --start-minimized
Icon=vocalinux
Comment=Voice dictation for Linux
Terminal=false
StartupNotify=false
X-GNOME-Autostart-enabled=true
